# supervisord.example.conf
[unix_http_server]
file=/var/run/supervisord/supervisord.sock

[supervisord]
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisord/supervisord.sock

[group:Bongos]
programs=web,application

[program:web]
command=python /path/to/Bongos/web.py

[program:application]
command=rqworker Bongos --path /path/to/Bongos/ --quiet
